# systemd/goldspider.service
[Unit]
Description=Gold Spider FastAPI Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/gold_spider
Environment=PATH=/home/ubuntu/anaconda3/envs/gold_spider/bin
ExecStart=/bin/bash -c 'source /home/ubuntu/anaconda3/etc/profile.d/conda.sh && conda activate gold_spider && python -m app.service'
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goldspider

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/gold_spider/logs /home/ubuntu/gold_spider/data

[Install]
WantedBy=multi-user.target 