# systemd/goldspider.service
[Unit]
Description=Gold Spider FastAPI Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Gold-Scraper
Environment=PATH=/root/anaconda3/envs/gold_spider/bin
ExecStart=/bin/bash -c 'source /root/anaconda3/etc/profile.d/conda.sh && conda activate gold_spider && python -m app.service'
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goldspider

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/Gold-Scraper/logs /root/Gold-Scraper/data

[Install]
WantedBy=multi-user.target 